# cuckoo-playbook/roles/start_up_cuckoo/files/cuckoo_api.service
[Unit]
Description=Cuckoo API Service
After=network.target cuckoo_web.service

[Service]
ExecStart=/opt/cuckoo/bin/cuckoo api --host 172.16.1.30 --port 8881
Restart=on-failure
User=cuckoo
Group=cuckoo
WorkingDirectory=/home/cuckoo/.cuckoo

[Install]
WantedBy=multi-user.target
