# cuckoo-playbook/roles/start_up_cuckoo/files/cuckoo_rooter.service
[Unit]
Description=Cuckoo Rooter Service
After=network.target openvpnclient.service

[Service]
ExecStart=/opt/cuckoo/bin/cuckoo rooter
Restart=on-failure
User=root
Group=root
WorkingDirectory=/home/cuckoo/.cuckoo

[Install]
WantedBy=multi-user.target
