# cuckoo-playbook/roles/openvpn/files/openvpnclient.service
[Unit]
Description=Openvpn Client Service
After=network.target

[Service]
ExecStart=/usr/sbin/openvpn --config openvpn.conf --script-security 2 --route-noexec --route-up route.py
Restart=on-failure
User=root
Group=root
WorkingDirectory=/etc/openvpn/client

[Install]
WantedBy=multi-user.target
