# cuckoo-playbook/roles/moloch/files/molochviewer.service
[Unit]
Description=Moloch Viewer
After=network.target elasticsearch.service
Requires=network.target elasticsearch.service

[Service]
Restart=on-failure
StandardOutput=tty
EnvironmentFile=-/data/moloch/etc/molochviewer.env
ExecStart=/bin/sh -c 'sleep 5 && /data/moloch/bin/node viewer.js -c /data/moloch/etc/config.ini ${OPTIONS} >> /data/moloch/logs/viewer.log 2>&1'
WorkingDirectory=/data/moloch/viewer

[Install]
WantedBy=multi-user.target
